# distiller-wifi.service
[Unit]
Description=Distiller WiFi Setup Service
After=network.target NetworkManager.service
Wants=network.target
Wants=NetworkManager.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/distiller-cm5-services
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/distiller-wifi-setup
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# Environment
Environment=WIFI_HOTSPOT_SSID=SetupWiFi
Environment=WIFI_HOTSPOT_PASSWORD=setupwifi123

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log /opt/distiller-cm5-services /etc/NetworkManager/system-connections

[Install]
WantedBy=multi-user.target 
