# pinggy-tunnel.service
[Unit]
Description=Pinggy SSH Tunnel Service
After=network-online.target wifi-setup.service
Wants=network-online.target
# Only start if network is actually connected
ConditionPathExists=/sys/class/net/wlan0/carrier

[Service]
Type=simple
User=root
WorkingDirectory=/home/distiller/distiller-cm5-services
ExecStart=/usr/bin/python3 /home/distiller/distiller-cm5-services/pinggy_tunnel_service.py --port 3000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/distiller/distiller-cm5-services"

# Give WiFi setup time to complete
ExecStartPre=/bin/sleep 30

[Install]
WantedBy=multi-user.target