# pinggy-tunnel.service
[Unit]
Description=Pinggy SSH Tunnel Service
After=network-online.target distiller-wifi.service
Wants=network-online.target
# Only start if network is actually connected
ConditionPathExists=/sys/class/net/wlan0/carrier

[Service]
Type=simple
User=distiller
Group=distiller
WorkingDirectory=/opt/distiller-cm5-services
ExecStart=/opt/distiller-cm5-services/.venv/bin/python pinggy_tunnel_service.py --port 3000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/distiller-cm5-sdk:/opt/distiller-cm5-sdk/src
Environment=LD_LIBRARY_PATH=/opt/distiller-cm5-sdk/lib

# Give WiFi setup time to complete
ExecStartPre=/bin/sleep 30

# Security settings
SupplementaryGroups=netdev gpio spi dialout audio

[Install]
WantedBy=multi-user.target
