# wifi-setup.service
[Unit]
Description=WiFi Setup Service
After=network-online.target
Wants=network-online.target
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/distiller/wifi_server
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 wifi_setup_service.py
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# Environment
Environment=PYTHONPATH=/home/distiller/wifi_server
Environment=WIFI_HOTSPOT_SSID=SetupWiFi
Environment=WIFI_HOTSPOT_PASSWORD=setupwifi123

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-setup

# Security (adjust as needed)
NoNewPrivileges=false
ProtectHome=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target 