# debian/postinst
#!/bin/sh
set -e

export PATH="/root/.local/bin:/root/.cargo/bin:$HOME/.local/bin:$HOME/.cargo/bin:/usr/local/bin:$PATH"

case "$1" in
configure)
	# Set MAC-based hostname only on target device
	if [ -f /boot/armbianEnv.txt ] || [ -f /boot/firmware/config.txt ]; then
		echo "Setting MAC-based hostname..."
		if [ -x /opt/distiller-cm5-services/scripts/update-hostname.sh ]; then
			/opt/distiller-cm5-services/scripts/update-hostname.sh || true
		fi
	fi

	cd /opt/distiller-cm5-services

	# Setup virtual environment
	[ -d ".venv" ] && rm -rf ".venv"
	uv venv --system-site-packages 2>/dev/null || uv venv

	# Install dependencies
	uv sync

	VENV_PATH="/opt/distiller-cm5-services/.venv"

	# Install SDK if available
	if [ -d "/opt/distiller-cm5-sdk" ]; then
		uv pip install -e /opt/distiller-cm5-sdk 2>/dev/null || true
	fi

	# Set permissions
	find /opt/distiller-cm5-services -type d -exec chmod 755 {} \; \
		-o -type f -name "*.py" -exec chmod 755 {} \;
	chmod +x /opt/distiller-cm5-services/distiller_wifi.py

	# Update systemd service
	if [ -f /lib/systemd/system/distiller-wifi.service ]; then
		sed -i "s|/usr/bin/python3|$VENV_PATH/bin/python|g" /lib/systemd/system/distiller-wifi.service
	fi

	# Reload systemd and enable service
	systemctl daemon-reload 2>/dev/null || true
	systemctl enable distiller-wifi.service 2>/dev/null || true

	echo "Distiller WiFi service installed. Start with: systemctl start distiller-wifi"
	;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

*)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

#DEBHELPER#

exit 0
