# debian/postinst
#!/bin/sh
set -e

export PATH="/root/.local/bin:/root/.cargo/bin:$HOME/.local/bin:$HOME/.cargo/bin:/usr/local/bin:$PATH"

case "$1" in
configure)
	# Create system directories
	install -d -o root -g root -m 755 /var/lib/distiller
	install -d -o root -g root -m 755 /var/log/distiller
	install -d -o root -g root -m 755 /etc/distiller

	cd /opt/distiller-cm5-services

	# Verify uv is available
	if ! command -v uv >/dev/null 2>&1; then
		echo "ERROR: uv not found in PATH" >&2
		exit 1
	fi

	# Setup virtual environment
	[ -d ".venv" ] && rm -rf ".venv"
	uv venv --system-site-packages 2>/dev/null || uv venv

	# Install dependencies
	if [ -f pyproject.toml ]; then
		uv sync || uv pip install -r requirements.txt
	else
		uv pip install -r requirements.txt
	fi

	VENV_PATH="/opt/distiller-cm5-services/.venv"

	# Install SDK if available
	if [ -d "/opt/distiller-cm5-sdk" ]; then
		uv pip install -e /opt/distiller-cm5-sdk 2>/dev/null || true
	fi

	# Set permissions
	find /opt/distiller-cm5-services -type d -exec chmod 755 {} \; \
		-o -type f -name "*.py" -exec chmod 755 {} \;
	chmod +x /opt/distiller-cm5-services/distiller_wifi.py

	# Update systemd service
	if [ -f /lib/systemd/system/distiller-wifi.service ]; then
		sed -i "s|/usr/bin/python3|$VENV_PATH/bin/python|g" /lib/systemd/system/distiller-wifi.service
	fi

	# Reload systemd and enable service
	systemctl daemon-reload 2>/dev/null || true
	systemctl enable distiller-wifi.service 2>/dev/null || true

	echo "Distiller WiFi service installed. Start with: systemctl start distiller-wifi"
	;;

abort-upgrade | abort-remove | abort-deconfigure)
	;;

*)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

#DEBHELPER#

exit 0