# debian/postrm
#!/bin/sh
set -e

case "$1" in
purge | remove)
	# Stop and disable service
	if systemctl is-active --quiet distiller-wifi.service 2>/dev/null; then
		systemctl stop distiller-wifi.service 2>/dev/null || true
	fi
	systemctl disable distiller-wifi.service 2>/dev/null || true

	# Remove installation
	[ -d "/opt/distiller-services" ] && rm -rf /opt/distiller-services
	[ -d "/opt/distiller-cm5-services" ] && rm -rf /opt/distiller-cm5-services

	if [ "$1" = "purge" ]; then
		echo "INFO: Purging user data..."
		# Clean up user data directories
		[ -d /etc/distiller ] && rm -rf /etc/distiller
		[ -d /var/lib/distiller ] && rm -rf /var/lib/distiller
		[ -d /var/log/distiller ] && rm -rf /var/log/distiller
		[ -f /etc/avahi/services/distiller-wifi.service ] &&
			rm -f /etc/avahi/services/distiller-wifi.service

		# Remove configuration files
		[ -f /etc/sudoers.d/distiller-system ] && rm -f /etc/sudoers.d/distiller-system
		[ -f /etc/polkit-1/rules.d/50-distiller-networkmanager.rules ] && rm -f /etc/polkit-1/rules.d/50-distiller-networkmanager.rules
		[ -f /etc/udev/rules.d/99-distiller-hardware.rules ] && rm -f /etc/udev/rules.d/99-distiller-hardware.rules

		# Clean up temporary files
		rm -f /tmp/distiller-* /var/run/distiller-* /var/lock/distiller-*

		# Remove NetworkManager connections
		if command -v nmcli >/dev/null 2>&1; then
			nmcli connection show 2>/dev/null | grep -E "^Distiller-" | awk '{print $1}' |
				xargs -r nmcli connection delete 2>/dev/null || true
		fi
	fi

	systemctl daemon-reload 2>/dev/null || true
	;;

upgrade | failed-upgrade | abort-install | abort-upgrade | disappear)
	# Stop service during upgrade but don't delete files
	if systemctl is-active --quiet distiller-wifi.service 2>/dev/null; then
		systemctl stop distiller-wifi.service 2>/dev/null || true
	fi
	systemctl daemon-reload 2>/dev/null || true
	;;

*)
	echo "ERROR: postrm called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

#DEBHELPER#

exit 0
